# silverblue-base/files/config.d/config.d.timer
[Unit]
Description=Run config.d.service hourly

[Timer]
OnBootSec=15min
OnCalendar=hourly

[Install]
WantedBy=multi-user.target